# deploy/systemd/clawdesk.service
[Unit]
Description=ClawDesk
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/clawdesk
EnvironmentFile=/etc/clawdesk/clawdesk.env
Environment=NODE_ENV=production
ExecStart=/usr/bin/npm run start -- --hostname 127.0.0.1 --port 3010
Restart=always
RestartSec=5
User=clawdesk
Group=clawdesk

[Install]
WantedBy=multi-user.target
